# silly-exchange-bot.service
[Unit]
Description=Silly Exchange Discord Bot
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/silly-exchange-bot
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=silly-exchange-bot

[Install]
WantedBy=multi-user.target 